# root/etc/services.d/libvirt-api/run
#!/usr/bin/with-contenv bash

echo "Starting flask Server."

[ -z "$API_KEY" ] && API_KEY="abcdef123456"

cd /app/

exec \
  gunicorn --workers 4 -b 0.0.0.0:80 hello:hello
