# root/etc/services.d/libvirt-api/run
#!/usr/bin/with-contenv bash

echo "Starting flask Server."

[ -z "$PORT" ] && PORT="80"

cd /app/libvirt-api

exec \
  gunicorn --workers 4 -b 0.0.0.0:$PORT wsgi:app
