# root/etc/services.d/libvirt-api/run
#!/usr/bin/with-contenv bash

echo "Starting flask Server."

[ -z "$API_KEY" ] && API_KEY="abcdef123456"

exec \
  /bin/bash /app/run.sh
