# root/etc/services.d/libvirt-mqtt/run
#!/usr/bin/with-contenv bash

echo "Starting MQTT service Server."

cd /app/

exec \
  python3 mqtt.py
